# docker/etc/cron.d/crontest
* * * * * root (date) >> /var/www/default/.crontest.log
