# docker/etc/cron.d/logrotate
* * * * * root (date; /usr/sbin/logrotate -v /etc/logrotate.conf 2>&1) >> /var/log/cron_logrotate.log

