# docker/etc/cron.d/logrotate
* * * * * root (date; logrotate -v /etc/logrotate.conf) >> /var/log/cron_logrotate.log
